# infra/ansible/roles/sibyl/templates/sibyl.service.j2
# {{ ansible_managed }}
[Unit]
Description=Sibyl knowledge graph stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory={{ sibyl_dir }}
ExecStart=/usr/bin/docker compose --env-file {{ sibyl_dir }}/.env up -d --build
ExecStop=/usr/bin/docker compose --env-file {{ sibyl_dir }}/.env down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
